# spexop-api.service
[Unit]
Description=Spexop Next.js API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/spexop-api
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=spexop-api

[Install]
WantedBy=multi-user.target